# mariadb/default/utils/supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /dev/shm/supervisor.sock

[supervisorctl]
serverurl = unix:///dev/shm/supervisor.sock
prompt = mariadb

[supervisord]
nodaemon=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mariadb]
command=/app/start_mysql.sh
process_name=mariadb
numprocs=1
priority=999
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
stopwaitsecs=300
stopasgroup=true
killasgroup=true

[program:listen]
command=/app/listen.sh
process_name=listen
numprocs=1
priority=999
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
startsecs=0

[program:monitor]
command=/app/mysql_monitor.sh
process_name=monitor
numprocs=1
priority=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
startsecs=0
startretries=99999

[program:filebeat]
command=/usr/bin/filebeat.sh -e -c /etc/filebeat.yml
process_name=filebeat
numprocs=1
priority=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
startsecs=0
